# redis.ini
[supervisord]
nodaemon=false

[program:redis-master-1]
command=redis-server /src/redis-master-1/redis.conf
process_name=%(program_name)s
directory=/src/redis-master-1
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/src/redis-master-1/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

[program:redis-master-2]
command=redis-server /src/redis-master-2/redis.conf
process_name=%(program_name)s
directory=/src/redis-master-2
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/src/redis-master-2/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

[program:redis-master-3]
command=redis-server /src/redis-master-3/redis.conf
process_name=%(program_name)s
directory=/src/redis-master-3
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/src/redis-master-3/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

[program:redis-replica-1]
command=redis-server /src/redis-replica-1/redis.conf
process_name=%(program_name)s
directory=/src/redis-replica-1
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/src/redis-replica-1/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

[program:redis-replica-2]
command=redis-server /src/redis-replica-2/redis.conf
process_name=%(program_name)s
directory=/src/redis-replica-2
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/src/redis-replica-2/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB

[program:redis-replica-3]
command=redis-server /src/redis-replica-3/redis.conf
process_name=%(program_name)s
directory=/src/redis-replica-3
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/src/redis-replica-3/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=50MB
